# filesystem/static/root/etc/profile.d/motd.sh
echo
echo "version: $(cat /etc/allium-release)"
echo
