# filesystem/static/root/etc/profile.d/motd.sh
# shellcheck disable=SC2148
echo
echo "version: $(cat /etc/pinewall-release)"
echo
