# filesystem/static/root/etc/profile.d/timezone.sh
# shellcheck disable=SC2148
export TZ='UTC'
