# config/vendor/etc/periodic/daily/pinehole
#!/usr/bin/env ash
set -eufo pipefail

# /etc/periodic/daily/pinehole

# This script runs daily at 2AM based on the Alpine Linux defaults
# See the /etc/periodic/daily entry in `crontab -l` for specifics

# Create tempfile and trap to remove it on exit
TMP_FILE="$(mktemp)"
trap 'rm -f "$TMP_FILE"' EXIT

# Download StevenBlack hosts file and process it into Unbound format
#
# We use NXDOMAIN here so that we sinkhole any domains underneath an
# ad domain that's on the list, rather than just serving 0.0.0.0 as an
# A record which might not be quite so robust
curl -fsSL https://raw.githubusercontent.com/StevenBlack/hosts/master/hosts | \
  awk '$1=="0.0.0.0" { print "local-zone: \"" $2 "\" always_nxdomain" }' > "$TMP_FILE"

# Check to ensure that the file has more than 10,000 entries
# and copy it into place if it does
if [ "$(wc -l ${TMP_FILE} | cut -d " " -f1)" -gt "10000" ]; then
    # Ensure file is readable by unprivileged Unbound process
    chmod 0644 "$TMP_FILE"
    # Move file into place over the top of the existing one
    mv -f "$TMP_FILE" /etc/unbound/adblock.list
    # Reload Unbound while keeping the cache
    unbound-control reload_keep_cache
else
    echo "List count is <10,000. Maybe something is broken with the download?"
    echo "Not loading adblock list into Unbound, just in case."
fi
